# src/etc/rc.d/kixdns
#!/bin/sh

# PROVIDE: kixdns
# REQUIRE: LOGIN networking
# KEYWORD: shutdown

. /etc/rc.subr

name="kixdns"
rcvar="kixdns_enable"

load_rc_config $name

: ${kixdns_enable:="NO"}
: ${kixdns_config:="/usr/local/etc/kixdns/pipeline.json"}
: ${kixdns_listener_label:="default"}
: ${kixdns_debug:="NO"}
: ${kixdns_udp_workers:="0"}

pidfile="/var/run/${name}.pid"
procname="/usr/local/bin/kixdns"

# Build command arguments
build_args() {
    _args="--config ${kixdns_config}"
    _args="${_args} --listener-label ${kixdns_listener_label}"
    
    if checkyesno kixdns_debug; then
        _args="${_args} --debug"
    fi
    
    if [ "${kixdns_udp_workers}" != "0" ]; then
        _args="${_args} --udp-workers ${kixdns_udp_workers}"
    fi
    
    echo "${_args}"
}

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"

kixdns_start() {
    if [ -f "${pidfile}" ]; then
        _pid=$(cat "${pidfile}" 2>/dev/null)
        if [ -n "${_pid}" ] && kill -0 "${_pid}" 2>/dev/null; then
            echo "${name} is already running (pid ${_pid})."
            return 1
        fi
        rm -f "${pidfile}"
    fi
    
    _args=$(build_args)
    echo "Starting ${name}."
    /usr/sbin/daemon -c -f -p "${pidfile}" -t "${name}" ${procname} ${_args}
}

kixdns_stop() {
    if [ -f "${pidfile}" ]; then
        _pid=$(cat "${pidfile}" 2>/dev/null)
        if [ -n "${_pid}" ]; then
            echo "Stopping ${name}."
            kill "${_pid}" 2>/dev/null
            # Wait for process to terminate
            _count=0
            while kill -0 "${_pid}" 2>/dev/null && [ ${_count} -lt 30 ]; do
                sleep 1
                _count=$((_count + 1))
            done
            if kill -0 "${_pid}" 2>/dev/null; then
                echo "Forcing ${name} to stop."
                kill -9 "${_pid}" 2>/dev/null
            fi
            rm -f "${pidfile}"
            return 0
        fi
    fi
    echo "${name} is not running."
    return 1
}

kixdns_status() {
    if [ -f "${pidfile}" ]; then
        _pid=$(cat "${pidfile}" 2>/dev/null)
        if [ -n "${_pid}" ] && kill -0 "${_pid}" 2>/dev/null; then
            echo "${name} is running (pid ${_pid})."
            return 0
        fi
    fi
    echo "${name} is not running."
    return 1
}

run_rc_command "$1"
